# data_popularity_api/deb/data-popularity-api/supervisor_conf/data-popularity-api.conf
[program:data-popularity-api]
command = /usr/local/bin/data-popularity-api
environment = DATA_POPULARITY_API_CONFIG="/etc/data-popularity-api.app.config"
user = robot-efweb
